# recipes-graphics/weston/files/weston.socket
[Unit]
Description=Weston Activation Socket

[Socket]
ListenStream=/run/user/root/wayland-0

[Install]
WantedBy=sockets.target
